# contrib/launchd/dev.localhttp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.localhttp</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/localhttp</string>
    <string>serve</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/localhttp.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/localhttp.launchd.err.log</string>
</dict>
</plist>
